# conf/collaborate-supervisor-dummy.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile          = ./logs/redaxo-collaborate-supervisord.log
logfile_maxbytes = 20MB
logfile_backups  = 5
loglevel         = info
pidfile          = /tmp/supervisord.pid
nodaemon         = false
minfds           = 1024
minprocs         = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:ratchet]
command                 = bash -c "ulimit -n 10000; exec /usr/bin/php ##COLLABORATE_INIT_SCRIPT_PATH##"
process_name            = Ratchet
numprocs                = 1
autostart               = true
autorestart             = true
user                    = root
stdout_logfile          = ./logs/redaxo-collaborate-ratchet-info.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = ./logs/redaxo-collaborate-ratchet-error.log
stderr_logfile_maxbytes = 1MB